# deploy/tradingbot.service
[Unit]
Description=Q-Trader HFT Bot — Microstructure Trading Engine
Documentation=https://github.com/your-repo/q-trader
After=network-online.target
Wants=network-online.target
# Disable restart rate limiting — allow infinite restarts
StartLimitIntervalSec=0

[Service]
# === Process Management ===
Type=notify
ExecStart=/opt/tradingbot/venv/bin/python /opt/tradingbot/run_bot.py
WorkingDirectory=/opt/tradingbot
User=tradingbot
Group=tradingbot

# === Watchdog ===
# If the bot doesn't send WATCHDOG=1 within 60s, systemd kills and restarts it.
# The bot sends pings every 20s via sdnotify.
WatchdogSec=60s

# === Restart Policy ===
Restart=always
RestartSec=5s

# === Resource Limits ===
MemoryMax=512M
CPUQuota=80%

# === Environment ===
EnvironmentFile=/opt/tradingbot/.env

# === Security Hardening ===
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/tradingbot/data
ReadWritePaths=/opt/tradingbot/.env

# === Logging ===
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tradingbot

[Install]
WantedBy=multi-user.target
